# doc/supervisor.ini
[program:ipa_server]
process_name=ipa_server-%(process_num)s
directory=/home/work/online/ipa/ipa_server/
command=python3 ipa_server.py --port=%(process_num)s --env=online --log_level=info --log_dir=/home/work/online/ipa/log
numprocs=6
numprocs_start=9100
user=work
startsecs=1
autorestart=true
autostart=true
startretries=5
redirect_stderr=true
stdout_logfile=/home/work/online/ipa/log/stdout_%(process_num)s.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=1
stderr_logfile=/home/work/online/ipa/log/stderr_%(process_num)s.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=1
